# gandi-dyndns.cron
*/5 * * * * gandi-dyndns -d <domain>.<tld> -r <record>
